# rootfs/etc/s6-overlay/s6-rc.d/cross-seed/run
#!/usr/bin/with-contenv sh

# run (from /usr/local/bin)
exec cross-seed daemon --verbose
