# rootfs/etc/s6-overlay/s6-rc.d/cross-seed/run
#!/usr/bin/with-contenv sh

# change chdir
cd /config

# run with limited permissions
exec s6-setuidgid $PUID:$PGID node /app/dist/cmd.js daemon --verbose
